# risved.service
[Unit]
Description=Risved Control Plane
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/risved
EnvironmentFile=/opt/risved/.env
ExecStart=/usr/bin/node /opt/risved/build
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
